# sparse/usr/lib/systemd/system/droid-late-start.service.d/10-droid-nagara-hal-init.conf
[Unit]
Requires=droid-nagara-hal-init.service

